# deploy/chatterbot.service
[Unit]
Description=ChatterBot Pi services (head + camera)
After=network-online.target zenohd.service
Wants=network-online.target
# Start after the router, but a router restart should not tear this down
# (the services reconnect on their own), so this is Wants, not Requires.
Wants=zenohd.service

[Service]
Type=simple
User=bruce
Group=bruce
WorkingDirectory=/home/bruce/Documents/Projects/ChatterBot
# Runs from the project venv, which must be created with --system-site-packages
# so the apt-installed picamera2 is visible:
#   python3 -m venv --system-site-packages .venv
#   .venv/bin/pip install -r requirements-pi.txt
# If you instead installed the deps into the system Python, use
# /usr/bin/python3 here.
ExecStart=/home/bruce/Documents/Projects/ChatterBot/.venv/bin/python -m chatterbot.launcher
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
